# supervisord.conf
; ========== 1. supervisord 主进程全局配置 ==========
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
autoshutdown=true
stopwaitsecs=10

; ========= 2. 本地unix socket通信 ============
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:root

; =============== 3. RPC 接口=============
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ========= 4. 客户端 supervisorctl 连接配置 =============
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


; Xvfb 虚拟显示器配置
[program:xvfb]
command=bash -c "rm -rf /tmp/.X1-lock && Xvfb :1 -screen 0 1280x1080x24"
autostart=true
autorestart=true
startsecs=2

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

priority=10

; Chrome 浏览器配置
[program:chrome]
command=chromium \
    --display=:1 \
    --window-size=1280,1080 \
    --start-maximized \
    --no-sandbox \
    --disable-dev-shm-usage \
    --disable-setuid-sandbox \
    --disable-accelerated-2d-canvas \
    --disable-gpu \
    --disable-features=WelcomeExperience,SigninPromo \
    --no-first-run \
    --no-default-browser-check \
    --disable-infobars \
    --test-type \
    --disable-popup-blocking \
    --disable-gpu-sandbox \
    --no-xshm \
    --new-window=false \
    --disable-notifications \
    --disable-extensions \
    --disable-component-extensions-with-background-pages \
    --disable-prompt-on-repost \
    --disable-dialogs \
    --disable-modal-dialogs \
    --disable-web-security \
    --disable-site-isolation-trials \
    --remote-debugging-address=0.0.0.0 \
    --remote-debugging-port=8222 %(ENV_CHROME_ARGS)s
autostart=true
autorestart=true
startretries=3
startsecs=10

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 内存限制
memlimit=2048MB

environment=DISPLAY=:1

priority=20

; socat 进程配置
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
startsecs=2

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

priority=30

; x11vnc
[program:x11vnc]
command=x11vnc -display :1 -nopw -shared -listen 0.0.0.0 -xkb -forever -rfbport 5900
autostart=true
autorestart=true
startsecs=3

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

priority=40

; websockify 配置，将 VNC 转换为 websocket 协议
[program:websockify]
command=websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
startsecs=3

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

priority=50

; FastAPI
[program:app]
command=uvicorn app.main:app --host 0.0.0.0 --port 8080 %(ENV_UVI_ARGS)s
; 进程工作目录
directory=/sandbox
user=root
autostart=true        ; supervisord 启动时自动拉起 fastapi
autorestart=true      ; 进程崩溃、被杀后自动重启
startretries=5        ; 启动失败最多重试5次
startsecs=3           ; 运行3秒才算启动成功
stopwaitsecs=8        ; 优雅关闭等待时间，超时强制 kill

; 进程标准输出/错误日志
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 完整回收 uvicorn 子进程,沙箱必备
stopsignal=SIGINT
stopasgroup=true
killasgroup=true

;环境变量
environment=HOME=/home/ubuntu

;启动优先级
priority=60


; 使用分组来统一启动 or 关闭进程
[group:services]
programs=xvfb,chrome,socat,x11vnc,websockify
